# bhulekh.service
[Unit]
Description=Bhulekh Village Scraper
After=bhulekh-boot.service network-online.target
Requires=bhulekh-boot.service
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=ubuntu
Environment=PATH=/home/ubuntu/.local/bin:/home/ubuntu/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=/home/ubuntu/justokay
ExecStart=/home/ubuntu/.local/bin/uv run python run_village_workers.py \
    --workers 28 \
    --db /home/ubuntu/justokay/work_queue.db \
    --data-dir /home/ubuntu/justokay/bhulekh_data \
    --headless \
    --fast \
    --districts 1 2 3 4 5 6 7 8 9 10 11 12 13 14 15 16 17 19 20 21 22 23 24 25 26 27 28 29 30
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bhulekh
MemoryMax=24G

[Install]
WantedBy=multi-user.target
